# debian/vaultwarden.service
# Please consider only editing this file through systemd drop-in files
# (e.g. `systemctl edit vaultwarden`). This will enable you to get updates
# to this file automatically, without your overrides being lost.

# Refer to the manpage `systemd.unit` for detailed documentation on each option.
[Unit]
Description=Vaultwarden API server
Documentation=https://github.com/dani-garcia/vaultwarden/wiki
# The service requires network interfaces to be configured.
After=network.target
# If you want to expose the server to the network directly, you may find the following useful:
#After=network-online.target

# The default database backend, SQLite, does not require any external services.
# However, other backends (MariaDB, MySQL, PostgreSQL) - if running on this system - should specify dependencies,
# such that Vaultwarden is only started after (`After=`) the DB has successfully (`Requires=`) started.
# Example for MariaDB:
#After=mariadb.service
#Requires=mariadb.service

# If running Vaultwarden behind a reverse proxy - running on this system - consider ordering it before also:
#Before=nginx.service

# Refer to the manpages `systemd.service` and `systemd.exec` for detailed documentation on each option.
[Service]
## Service execution environment
# Service configuration file.
EnvironmentFile=-%E/default/vaultwarden
# Vaultwarden configuration file.
Environment=ENV_FILE=%E/vaultwarden/vaultwarden.env
# Service execution.
Type=exec
ExecStart=/usr/bin/vaultwarden
WorkingDirectory=/var/lib/vaultwarden

## Permissions and directories
User=vaultwarden
Group=vaultwarden
# Read-only configuration directory (read-only mount, permissions managed by tmpfiles, warns if mode differs).
ConfigurationDirectory=vaultwarden::ro
ConfigurationDirectoryMode=0500
# Writable directories (read-write for the service, read-only for the group, reapplies permissions recursively if different).
StateDirectory=vaultwarden
StateDirectoryMode=0750
CacheDirectory=vaultwarden
CacheDirectoryMode=0750
# File creation mask (matching modes above).
UMask=0027

## Connection and process limits
LimitNOFILE=1048576
# Systems with many cores may want to raise this limit:
# https://github.com/dionysius/vaultwarden-deb/issues/21
LimitNPROC=64
#LimitLOCKS=64
CPUQuota=200%
MemoryHigh=2G
MemoryMax=2.5G
#TasksMax=512

## Sandboxing and isolation
# An overview is available using `systemd-analyze security vaultwarden`.

## Capabilities
AmbientCapabilities=
CapabilityBoundingSet=
# Uncomment these to allow vaultwarden to bind to ports in the range 0-1024.
#AmbientCapabilities=CAP_NET_BIND_SERVICE
#CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Process restrictions
NoNewPrivileges=true
RestrictSUIDSGID=true
SecureBits=noroot noroot-locked

# Filesystem protection
ProtectSystem=strict
ProtectHome=true
ProtectProc=noaccess
ProcSubset=pid

# Namespace isolation
PrivateTmp=disconnected
PrivateDevices=true
PrivateIPC=true
PrivateBPF=true

# Kernel and system protection
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=strict

# Execution restrictions
MemoryDenyWriteExecute=true
SetLoginEnvironment=false

# System call and network filtering
SystemCallFilter=@system-service
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictFileSystems=@basic-api @common-block
RestrictRealtime=true

# Cleanup
RemoveIPC=true

[Install]
WantedBy=multi-user.target
